# Wigig_release_2.0.2/work/wil6210/10-wil6210.rules
SUBSYSTEM=="net", ENV{SOURCE}=="wil6210", ENV{EVENT}=="FW_ERROR", RUN+="/home/wigig/work/wil6210/fw_error.sh"
